# startmusic.service
[Unit]
Description=cmus

[Service]
Type=forking
User=hypen9
ExecStart=/home/hypen9/music.sh
Environment="HOME=/home/hypen9" "USERNAME=hypen9"

[Install]
WantedBy=multi-user.target
